# postgres-xc-1.2/postgres-xc-1.2@.service
[Unit]
Description=Postgres-XC (%i) datanode server
After=network.target

[Service]
Type=forking
TimeoutSec=120
User=postgres
Group=postgres

Environment=PGVER=1.2
Environment=PGROOT=/var/lib/postgres/xc/${PGVER}/%i

SyslogIdentifier=pgxc-${PGVER}-%i
PIDFile=${PGROOT}/data/postmaster.pid

ExecStart= /opt/postgres-xc-${PGVER}/bin/pg_ctl -s -Z "datanode" -D ${PGROOT} start -w -t 120
ExecReload=/opt/postgres-xc-${PGVER}/bin/pg_ctl -s -D ${PGROOT} reload
ExecStop=  /opt/postgres-xc-${PGVER}/bin/pg_ctl -s -D ${PGROOT} stop -m fast

# Due to PostgreSQL's use of shared memory, OOM killer is often overzealous in
# killing Postgres, so adjust it downward
OOMScoreAdjust=-200

[Install]
WantedBy=multi-user.target
